# SysadminsBot/sysadmins.service
[Unit]
Description=Sysadmins chat bot
 
[Service]
Type=notify
ExecStart=/opt/sysadmins/SysadminsBot
 
[Install]
WantedBy=multi-user.target